# water_com.conf
[program:water_com]
command = gunicorn water_com.wsgi_prod -b 127.0.0.1:8013
directory = /root/water_com/water3
autostart = true
autorestart = true
