# water.conf
[program:water]
command = gunicorn water_com.wsgi_prod -b 127.0.0.1:8013
directory = /root/water_com/water_com
autostart = true
autorestart = true
